# rightsize/supervisord.agent.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:agent]
command=python /usr/src/app/agent/agent.py
directory=/usr/src/app
autorestart=false
